# scripts/automation/com.gnosis.report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gnosis.report</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/y.noguchi/.bun/bin/bun</string>
        <string>src/scripts/status-report.ts</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/y.noguchi/Code/gnosis</string>
    <key>StartInterval</key>
    <integer>86400</integer>
    <key>StandardOutPath</key>
    <string>/Users/y.noguchi/Code/gnosis/logs/report.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/y.noguchi/Code/gnosis/logs/report.log</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
