# scripts/automation/com.gnosis.embedding-batch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gnosis.embedding-batch</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{BUN_PATH}}</string>
        <string>run</string>
        <string>src/scripts/embedding-batch-worker.ts</string>
    </array>
    <key>WorkingDirectory</key>
    <string>{{PROJECT_ROOT}}</string>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>StandardOutPath</key>
    <string>{{PROJECT_ROOT}}/logs/embedding-batch.log</string>
    <key>StandardErrorPath</key>
    <string>{{PROJECT_ROOT}}/logs/embedding-batch.log</string>
    <key>RunAtLoad</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>GNOSIS_ENABLE_AUTOMATION</key>
        <string>true</string>
        <key>GNOSIS_EMBED_DAEMON_URL</key>
        <string>http://127.0.0.1:44512</string>
        <key>GNOSIS_EMBED_BATCH_SIZE</key>
        <string>50</string>
    </dict>
</dict>
</plist>
