# scripts/automation/com.gnosis.guidance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gnosis.guidance</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{BUN_PATH}}</string>
        <string>src/scripts/import-guidance.ts</string>
        <string>--max-zips</string>
        <string>1</string>
    </array>
    <key>WorkingDirectory</key>
    <string>{{PROJECT_ROOT}}</string>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>{{PROJECT_ROOT}}/logs/guidance-import.log</string>
    <key>StandardErrorPath</key>
    <string>{{PROJECT_ROOT}}/logs/guidance-import.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
