# scripts/automation/com.gnosis.embedding-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gnosis.embedding-daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{PROJECT_ROOT}}/services/embedding/.venv/bin/python</string>
        <string>-m</string>
        <string>e5embed.daemon</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>44512</string>
    </array>
    <key>WorkingDirectory</key>
    <string>{{PROJECT_ROOT}}</string>
    <key>StandardOutPath</key>
    <string>{{PROJECT_ROOT}}/logs/embedding-daemon.log</string>
    <key>StandardErrorPath</key>
    <string>{{PROJECT_ROOT}}/logs/embedding-daemon.log</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
