# scripts/automation/com.gnosis.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gnosis.sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/y.noguchi/.bun/bin/bun</string>
        <string>run</string>
        <string>src/scripts/sync.ts</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/y.noguchi/Code/gnosis</string>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>/Users/y.noguchi/Code/gnosis/logs/sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/y.noguchi/Code/gnosis/logs/sync.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
